# csgoemporiumfrontend/deploy/configs/supervisor/guardmail.conf
[program:guardmail]
command=/var/www/emporium/scripts/guardmail.py
autostart=true
autorestart=true
user=emporium
environment=ENV="PROD"
