# csgoemporiumfrontend/deploy/configs/supervisor/pgbouncer.conf
[program:pgbouncer]
command=pgbouncer /etc/pgbouncer/pgbouncer.ini
autostart=true
autorestart=true
user=emporium

