# csgoemporiumfrontend/deploy/configs/supervisor/queue.conf
[program:queue-%(process_num)s]
command=/var/www/emporium/app/run queue
numprocs=4
autostart=true
autorestart=true
user=emporium
environment=ENV="PROD"
process_name=queue-%(process_num)s
