# csgoemporiumfrontend/deploy/configs/supervisor/rush.conf
[program:rush]
command=/var/www/emporium/rush/rush.bin
autostart=true
autorestart=true
user=emporium
