# csgoemporiumfrontend/deploy/configs/supervisor/scheduler.conf
[program:scheduler]
command=/var/www/emporium/app/run sched
autostart=true
autorestart=true
user=emporium
environment=ENV="PROD"

