# csgoemporiumfrontend/deploy/configs/supervisor/uwsgi.conf
[program:uwsgi]
command=/usr/bin/pidproxy /var/run/emporium/uwsgi.pid /var/www/emporium/uwsgi.bin
autostart=true
autorestart=true
user=emporium
environment=ENV="PROD"

